# nomad/nomad.upstart.conf
# /etc/inid/nomad.conf

start on filesystem and static-network-up

exec nomad agent -config /etc/nomad
